# data/s6/frost-live/up
#!/usr/bin/execlineb -P
if { echo "Initializing Frost Linux live system" }
foreground { /usr/bin/frost-live }
